# whisper-dictate.service
[Unit]
Description=Whisper Dictation Daemon
Documentation=https://github.com/openai/whisper
After=graphical-session.target whisper-streaming-server.service
Requires=whisper-streaming-server.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.local/bin/whisper_dictate_daemon.py

# Pass Wayland environment for wtype to work
PassEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR DISPLAY

# Environment configuration
# Model options: tiny.en, base.en, small.en, medium.en, large
Environment="JNICKG_DICTATE_MODEL=base.en"
Environment="JNICKG_DICTATE_SOCKET=/tmp/jnickg-dictate.sock"
# Input method: wtype (Wayland), xdotool (X11), ydotool (Wayland alternative)
Environment="JNICKG_DICTATE_INPUT_METHOD=ydotool"
# Streaming server configuration
Environment="JNICKG_DICTATE_STREAMING_HOST=localhost"
Environment="JNICKG_DICTATE_STREAMING_PORT=43001"

Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=false
ReadWritePaths=/tmp

[Install]
WantedBy=graphical-session.target
